# server/symfony-worker.conf
[program:symfony-worker]
command=php bin/console app:server
directory=.
autostart=true
autorestart=true
user=root